# Line/scheduler/crontab
0 21 * * * python3 /app/scheduler.py >> /var/log/cron.log 2>&1